# configs/s6/nginx/finish
#!/bin/sh

kill `cat /var/run/nginx/nginx.pid`